# com.aomb.morning-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aomb.morning-report</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/abhipan/AOMB/.venv/bin/python</string>
        <string>/Users/abhipan/AOMB/morning_report.py</string>
        <string>--plot</string>
    </array>

    <!-- Run every day at 8:00 AM local time (PST/PDT) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/abhipan/AOMB</string>

    <key>StandardOutPath</key>
    <string>/Users/abhipan/AOMB/logs/morning_report.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/abhipan/AOMB/logs/morning_report.err</string>

    <!-- Run even after a missed schedule (e.g. Mac was asleep at 8 AM) -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
